# pc-hardware-monitor.service
[Unit]
Description=PC Hardware Monitor for ESP32
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=machetie
WorkingDirectory=/home/machetie/Documents/git/esp32-pc-hardware-monitor
ExecStart=/usr/bin/env python3 /home/machetie/Documents/git/esp32-pc-hardware-monitor/pc_monitor.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target

